# root/etc/services.d/haproxy/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/sbin/haproxy -p /var/run/haproxy/haproxy.pid -f /etc/haproxy/haproxy.cfg 
