# root/etc/services.d/haproxy/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
