# root/etc/services.d/webhook/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /app/webhook -ip 0.0.0.0 -port 9000 -hotreload -template -hooks /config/hooks/hooks.yaml -nopanic -verbose
